# root/etc/systemd/system/download-updates.timer
[Unit]
Description=download updates for install

[Timer]
OnCalendar=Mon..Fri 12:30

[Install]
WantedBy=timers.target
