# root/etc/systemd/system/download-updates.service
[Unit]
Description=download package updates

[Service]
Type=simple
ExecStart=/usr/bin/pacman -Syuw --noconfirm
